# openclaw/deploy/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway for Deepwork Research
After=network.target deepwork-daemon.service

[Service]
Type=simple
User=deepwork
WorkingDirectory=/home/deepwork/deepwork/openclaw
ExecStart=/usr/local/bin/openclaw gateway --config gateway.json
Restart=always
RestartSec=10
EnvironmentFile=/home/deepwork/deepwork/.env
Environment=NODE_OPTIONS=--max-old-space-size=256

[Install]
WantedBy=multi-user.target
