# forge/deploy/forge-site.service
[Unit]
Description=Forge Dashboard — Next.js
Documentation=https://github.com/serre-ai/research
After=network-online.target forge-daemon.service
Wants=network-online.target

[Service]
Type=simple
User=deepwork
Group=deepwork
WorkingDirectory=/opt/deepwork/site-next/.next/standalone/site-next
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
TimeoutStopSec=30

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
EnvironmentFile=/opt/deepwork/site-next/.env.local

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome=no because /opt/deepwork is a symlink to /home/deepwork/deepwork
# and ProtectHome=true blocks symlink resolution even with ReadWritePaths
ProtectHome=no
ReadWritePaths=/opt/deepwork /home/deepwork/deepwork
PrivateTmp=true

# Resource limits
MemoryMax=512M
MemoryHigh=384M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=forge-site

[Install]
WantedBy=multi-user.target
