# orchestrator/deploy/deepwork-daemon.service
[Unit]
Description=Deepwork Research Platform Daemon
Documentation=https://github.com/oddurs/deepwork
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=deepwork
Group=deepwork
WorkingDirectory=/opt/deepwork
ExecStart=/usr/bin/node orchestrator/dist/index.js run --root-dir /opt/deepwork-agent
Restart=always
RestartSec=30
TimeoutStopSec=120

# Environment
EnvironmentFile=/opt/deepwork/.env

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome=no because /opt/deepwork and /opt/deepwork-agent are symlinks
# into /home/deepwork/ which ProtectHome=true blocks
ProtectHome=no
ReadWritePaths=/opt/deepwork /opt/deepwork-agent /home/deepwork
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deepwork-daemon

[Install]
WantedBy=multi-user.target
