# orchestrator/deploy/deepwork-eval@.service
[Unit]
Description=Deepwork Eval Job — %i
Documentation=https://github.com/oddurs/deepwork
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=deepwork
Group=deepwork
WorkingDirectory=/opt/deepwork/projects/reasoning-gaps/benchmarks

# Instance name format: model::task::condition (using :: delimiter)
# Example: systemctl start deepwork-eval@openai:o3::B5::budget_cot
# Runs eval then auto-ingests results into PostgreSQL
ExecStart=/bin/bash -c '/opt/deepwork/.venv/bin/python3 run_evaluation.py --instance %i && /opt/deepwork/.venv/bin/python3 ingest_results.py'

Restart=no
RestartSec=60

# Allow graceful checkpoint save on stop
KillSignal=SIGINT
SendSIGHUP=yes
TimeoutStopSec=300

# Environment
EnvironmentFile=/opt/deepwork/.env

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/deepwork
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deepwork-eval-%i

[Install]
WantedBy=multi-user.target
